# root/etc/s6/nave-npm-start/run
#!/bin/bash
cd /app
echo "runned nave"
nave use stable npm start
